# root/etc/services.d/lsyncd/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid root lsyncd -nodaemon -insist /config/lsyncd.lua
